# overlay/usr/lib/systemd/system/kubelet.service.d/11-resource-sizing.conf
[Service]
EnvironmentFile=-/run/kubelet/extra-args.env
ExecStartPre=/usr/libexec/kubernetes/kubelet-resource-sizing.sh